# systemd/isli-skills.service
[Unit]
Description=ISLI Skills
After=network.target redis-server.service isli-core.service
Wants=redis-server.service

[Service]
Type=simple
User=isli
Group=isli
WorkingDirectory=/opt/isli/isli-skills
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/isli/.env
ExecStart=/opt/isli/isli-skills/.venv/bin/python -m uvicorn isli_skills.main:app --host 0.0.0.0 --port 8003
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
