# systemd/isli-core.service
[Unit]
Description=ISLI Core API
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=isli
Group=isli
WorkingDirectory=/opt/isli/isli-core
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/isli/.env
ExecStart=/opt/isli/isli-core/.venv/bin/python -m uvicorn isli_core.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
